# scheduler/crontab
# Run backup every day at 01:00
# m h dom mon dow command
0 1 * * * /app/backup_all_volumes.sh >> /var/log/cron.log 2>&1
