# scheduler/crontab
# Run backup every day at 23:00
# m h dom mon dow command
0 23 * * * /app/backup_all_volumes.sh >> /var/log/cron.log 2>&1
